# arch/install.sh
#!/bin/bash

yay -S --needed \
  adw-gtk-theme \
  android-tools \
  apple-fonts \
  base-devel \
  bat \
  bibata-cursor-theme \
  bluetui \
  bluez \
  bluez-utils \
  bottom \
  brave-bin \
  chafa \
  cmake \
  cowsay \
  cursor-bin \
  datagrip \
  datagrip-jre \
  dbus \
  deno \
  distrobox \
  dive \
  eglinfo \
  elephant \
  elephant-archlinuxpkgs \
  elephant-bluetooth \
  elephant-calc \
  elephant-clipboard \
  elephant-desktopapplications \
  elephant-files \
  elephant-menus \
  elephant-providerlist \
  elephant-runner \
  elephant-symbols \
  elephant-todo \
  elephant-unicode \
  elephant-websearch \
  elephant-windows \
  eza \
  fd \
  ffmpeg \
  ffmpegthumbnailer \
  fnm \
  fzf \
  gcc \
  gimp \
  github-cli \
  glxinfo \
  gnome-keyring \
  go \
  google-chrome \
  gparted \
  grim \
  gvfs \
  gzip \
  imv \
  inter-font \
  inxi \
  iw \
  jfsutils \
  jq \
  killport \
  krita \
  lazydocker \
  lazygit \
  less \
  lib32-vulkan-intel \
  libqalculate \
  libxcrypt-compat \
  local-by-flywheel-bin \
  lshw \
  lua \
  luarocks \
  lynx \
  mkcert \
  mpv \
  nemo \
  nemo-fileroller \
  neovim \
  neovim \
  networkmanager \
  nilfs-utils \
  notesnook \
  noto-fonts \
  noto-fonts-cjk \
  noto-fonts-emoji \
  noto-fonts-extra \
  ntfs-3g \
  ntfsprogs \
  obsidian \
  opencode-bin \
  openssh \
  p7zip \
  papipurs-icon-theme \
  pavucontrol \
  pfetch \
  podman-compose \
  polkit-gnome \
  postman \
  python-pip \
  python-pipx \
  qwen-code \
  reflector \
  reiserfsprogs \
  ripgrep \
  rustup \
  sed \
  slurp \
  sof-firmware \
  starship \
  stow \
  swaybg \
  swayfx \
  swayidle \
  swaylock-effects-git \
  swaync \
  swayosd-git \
  tela-icon-theme \
  termius \
  tmux \
  trash-cli \
  ttf-cascadia-code-nerd \
  ttf-firacode-nerd \
  udevil \
  udisks2 \
  unrar \
  unzip \
  visual-studio-code-bin \
  vivid \
  walker \
  waybar \
  wf-recorder \
  wget \
  wiremix \
  wireplumber \
  wl-clipboard \
  xdg-desktop-portal \
  xdg-desktop-portal-gtk \
  xdg-desktop-portal-wlr \
  xdg-user-dirs \
  xfsprogs \
  xorg-xhost \
  yaak-bin \
  yazi \
  youtube-dl \
  zathura \
  zathura-djvu \
  zathura-pdf-mupdf \
  zen-browser-bin \
  zip

sudo systemctl enable --now swayosd-libinput-backend.service
